# redis-cluster.service
[Unit]
Description=Redis Cluster Service
After=network.target

[Service]
ExecStart=/path/to/start-redis-cluster.sh
ExecStop=/usr/bin/killall redis-server
Restart=always
User=redis

[Install]
WantedBy=multi-user.target